# config/supervisor.conf
[program:litheblas]
directory = /opt/litheblas.org/
user = www-data
command = /opt/litheblas.org/bin/gunicorn -c /opt/litheblas.org/litheblas/config/gunicorn.py litheblas.wsgi
stdout_logfile = /opt/litheblas.org/log/supervisor/stdout.log
stderr_logfile = /opt/litheblas.org/log/supervisor/stderr.log
autostart = True
autorestart = True